# actix-share-listener-sample.socket
[Unit]
Description=Example App Socket

[Socket]
ListenStream=/tmp/share-fd.sock
SocketMode=0666

[Install]
WantedBy=sockets.target

